# config/includes.chroot/etc/systemd/system/atlazes-firstboot.service
[Unit]
Description=ATLAZES OS First Boot Setup
After=network.target
ConditionPathExists=!/var/lib/atlazes/.firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/atlazes-firstboot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
